# environments/t14s/notify/daily-notify.service
[Unit]
Description=Send daily notification

[Service]
Type=oneshot
ExecStart=%h/bin/daily-notify.sh
