# environments/t14s/notify/daily-notify.timer
[Unit]
Description=Run notify-send daily

[Timer]
OnCalendar=21:59
Persistent=true

[Install]
WantedBy=timers.target
